# deploy/systemd/eg-web.service
[Unit]
Description=Establecimiento Ganadero Web (Next.js)
After=network.target eg-api.service
Requires=eg-api.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/workspace/EstablecimientoGanadero
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/workspace/EstablecimientoGanadero/.env
ExecStart=/usr/bin/env npm --workspace apps/web run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
